# home/stream/daorakids-sync.timer
[Unit]
Description=Timer para Sincronizacao de Videos Daora Kids

[Timer]
# Roda 1 minuto apos o boot para garantir rede e libs
OnBootSec=1min
# Roda no minuto zero de cada hora
OnCalendar=hourly
# Garante que rode mesmo se o Pi ligar depois da hora marcada
Persistent=true
Unit=daorakids-sync.service

[Install]
WantedBy=timers.target
