# home/stream/daorakids-sync.timer
[Unit]
Description=Timer para Sincronizacao de Videos Daora Kids (Boot + Hora em Hora)

[Timer]
# Roda 2 minutos apos o boot (para dar tempo da internet subir)
OnBootSec=2min
# Roda exatamente no inicio de cada hora (Ex: 14:00, 15:00...)
OnCalendar=hourly
# Garante que rode mesmo se o Pi ligar depois da hora marcada
Persistent=true
Unit=daorakids-sync.service

[Install]
WantedBy=timers.target
